# server/supervisord.conf
[program:twitter_candidates]
command=/usr/bin/python /srv/www/kazanski/twitter.py
directory=/srv/www/kazanski
environment=PATH="/srv/www/kazanski_env/bin"
